# fsimg/etc/inittab
::sysinit:/etc/rcS
hvc0::respawn:-/bin/sh
::respawn:/sbin/getty 38400 tty1
::respawn:/sbin/getty 38400 tty2
::respawn:/sbin/getty 38400 tty3
::respawn:/sbin/getty 38400 tty4
::respawn:/sbin/getty 38400 tty5
